# deploy/systemd/linkflow.service
[Unit]
Description=LinkFlow full-stack app
After=network.target

[Service]
Type=simple
User=linkflow
Group=linkflow
WorkingDirectory=/var/www/linkflow/current
EnvironmentFile=/var/www/linkflow/shared/.env
ExecStart=/usr/bin/node server/index.mjs
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/www/linkflow/shared

[Install]
WantedBy=multi-user.target
